# systemd/shpool.socket
[Unit]
Description=Shpool Shell Session Pooler

[Socket]
ListenStream=%t/shpool/shpool.socket

[Install]
WantedBy=sockets.target
